# templates/etc/init.d/carbon-cache.erb
#!/bin/sh
# !!! generated by puppet !!!

### BEGIN INIT INFO
# chkconfig: 2345 99 80
# Provides: carbon-cache
# Required-Start: $remote_fs $network $time
# Required-Stop: $remote_fs $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Start/Stop carbon-cache
# Description: Enables Graphites carbon-cache data collecting engine
### END INIT INFO

<% 
	cache_list=""
	cc = 0
	while cc < scope.lookupvar('graphite::cache_count').to_i
		cache_list = "#{cache_list} #{cc}"
		cc += 1
	end
	cache_list[0]=""
%>

case "$1" in
	start)
		for cache in <%= cache_list %>; 
		do
			python /opt/graphite/bin/carbon-cache.py --instance=$cache start
		done
		;;
	stop)
		for cache in <%= cache_list %>; 
		do
			python /opt/graphite/bin/carbon-cache.py --instance=$cache stop
		done
		;;
	status)
		for cache in <%= cache_list %>; 
		do
			python /opt/graphite/bin/carbon-cache.py --instance=$cache status
		done
		;;
	restart)
		for cache in <%= cache_list %>; 
		do
			python /opt/graphite/bin/carbon-cache.py --instance=$cache stop
			sleep 3
			python /opt/graphite/bin/carbon-cache.py --instance=$cache  start
		done
		
		RC=$?
		exit $RC
		;;
	*)
		echo "Usage: $0 {start|stop|restart|status}"
		exit 1
		;;
esac

exit
