# templates/etc/init.d/carbon-relay.erb
#!/bin/sh
# !!! generated by puppet !!!

### BEGIN INIT INFO
# chkconfig: 2345 99 80
# Provides: carbon-relay
# Required-Start: $remote_fs $network $time
# Required-Stop: $remote_fs $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Start/Stop carbon-relay
# Description: Enables Graphites carbon-relay data routing engine
### END INIT INFO

<% 
	relay_list=""
	scope.lookupvar('graphite::relay_list').each do | name | 
	relay_list = "#{relay_list} #{name}"
end 
relay_list[0]=""
%>


case "$1" in
	start)
		for relay in <%= relay_list %>; 
		do
			python /opt/graphite/bin/carbon-relay.py --instance=$relay start
		done
		;;
	stop)
		for relay in <%= relay_list %>; 
		do
			python /opt/graphite/bin/carbon-relay.py --instance=$relay stop
		done
		;;
	status)
		for relay in <%= relay_list %>; 
		do
			python /opt/graphite/bin/carbon-relay.py --instance=$relay status
		done
		;;
	restart)
		for relay in <%= relay_list %>; 
		do
			python /opt/graphite/bin/carbon-relay.py --instance=$relay stop
			sleep 3
			python /opt/graphite/bin/carbon-relay.py --instance=$relay start
		done
		RC=$?
		exit $RC
		;;
	*)
		echo "Usage: $0 {start|stop|restart|status}"
		exit 1
		;;
esac

exit
