# templates/etc/init.d/carbon-aggregator.erb
#!/bin/sh
# !!! generated by puppet !!!

### BEGIN INIT INFO
# chkconfig: 2345 99 80
# Provides: carbon-aggregator
# Required-Start: $remote_fs $network $time
# Required-Stop: $remote_fs $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Start/Stop carbon-aggregator
# Description: Enables Graphites carbon-aggregator data collecting engine
### END INIT INFO

<% 
	aggregator_list=""
	cc = 0
	while cc < scope.lookupvar('graphite::aggregator_count').to_i
		aggregator_list = "#{aggregator_list} #{cc}"
		cc += 1
	end
	aggregator_list[0]=""
%>

case "$1" in
	start)
		for aggregator in <%= aggregator_list %>; 
		do
			python /opt/graphite/bin/carbon-aggregator.py --instance=$aggregator start
		done
		;;
	stop)
		for aggregator in <%= aggregator_list %>; 
		do
			python /opt/graphite/bin/carbon-aggregator.py --instance=$aggregator stop
		done
		;;
	status)
		for aggregator in <%= aggregator_list %>; 
		do
			python /opt/graphite/bin/carbon-aggregator.py --instance=$aggregator status
		done
		;;
	restart)
		for aggregator in <%= aggregator_list %>; 
		do
			python /opt/graphite/bin/carbon-aggregator.py --instance=$aggregator stop
			sleep 3
			python /opt/graphite/bin/carbon-aggregator.py --instance=$aggregator  start
		done
		
		RC=$?
		exit $RC
		;;
	*)
		echo "Usage: $0 {start|stop|restart|status}"
		exit 1
		;;
esac

exit
