# templates/etc/init.d/carbon-relay.erb
#!/bin/sh
# !!! generated by puppet !!!

### BEGIN INIT INFO
# chkconfig: 2345 99 80
# Provides: carbon-relay
# Required-Start: $remote_fs $network $time
# Required-Stop: $remote_fs $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Start/Stop carbon-relay
# Description: Enables Graphites carbon-relay data routing engine
### END INIT INFO

<% 
	start_collection = ""
	stop_collection = ""
	status_collection = ""
	scope.lookupvar('graphite::relays').each do | name, relay | 
	start_collection = "#{start_collection}\n\t\tpython /opt/graphite/bin/carbon-relay.py --instance=#{name} start"
	stop_collection = "#{stop_collection}\n\t\tpython /opt/graphite/bin/carbon-relay.py --instance=#{name} stop"
	status_collection = "#{status_collection}\n\t\tpython /opt/graphite/bin/carbon-relay.py --instance=#{name} status"
end %>

case "$1" in
	start)
		<%= start_collection %>
		;;
	stop)
		<%= stop_collection %>
		;;
	status)
		<%= status_collection %>
		;;
	restart)
		$0 stop ; sleep 3 ; $0 start
		RC=$?
		exit $RC
		;;
	*)
		echo "Usage: $0 {start|stop|restart|status}"
		exit 1
		;;
esac

exit
